# edge_deploy/systemd/mimii-api.service
[Unit]
Description=MIMII Edge HTTP Inference API
After=network.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/mimii
Environment=PYTHONUNBUFFERED=1
Environment=MIMII_ROOT=/opt/mimii
Environment=MIMII_CONFIG=/opt/mimii/config.yaml
Environment=OMP_NUM_THREADS=3
Environment=OPENBLAS_NUM_THREADS=3
ExecStart=/opt/mimii/venv/bin/uvicorn inference_server:app --host 0.0.0.0 --port 8000
User=mimii
Group=audio
Restart=always
RestartSec=5
MemoryMax=1G
CPUAffinity=1-3

[Install]
WantedBy=multi-user.target
